# settings/supervisor/create-topics.conf

[program:createtopics]
command=/bin/bash /root/create-topics.sh
autostart=true
startsecs=100
startretries=0
autorestart=false